# archive/forum-on-prem/deploy/forum-analysis.timer
[Unit]
Description=DISABLED — legacy forum-ai pipeline (edge cron replaces this)
Documentation=forum-airlock/docs/civic-edge-analysis.md

[Timer]
# Inactive schedule: public reports come from secure-worker every 6h.
# To re-enable legacy on-prem analysis, restore OnBootSec/OnUnitActiveSec and run:
#   systemctl --user enable --now forum-analysis.timer
OnCalendar=2099-01-01 00:00:00 UTC
Unit=forum-analysis.service

[Install]
WantedBy=timers.target
